# resumegen.service
[Unit]
Description=Gunicorn instance to serve Resume Generator AI
After=network.target

[Service]
# IMPORTANT: Replace 'your_server_user' with the actual user that should run the app
User=your_server_user
# IMPORTANT: Replace 'your_server_group' with the actual group for the user (often same as user)
Group=your_server_group
# IMPORTANT: Replace with the *absolute* path to your project directory on the server
WorkingDirectory=/path/to/your/Resume-Generator-AI
# IMPORTANT: Set environment variables needed by the app
# Ensure the .env file is in the WorkingDirectory or provide absolute path
Environment="PATH=/path/to/your/python/env/bin" # Add path to python/gunicorn if using venv
Environment="FLASK_APP=backend.app:create_app"
# Load other variables from .env file if needed, or set them here explicitly
# Environment="SECRET_KEY=your_strong_random_secret_key" # Example if not using .env
# Environment="GEMINI_API_KEY=YOUR_GOOGLE_GEMINI_API_KEY" # Example if not using .env

# IMPORTANT: Replace with the *absolute* path to the flask executable (often inside venv/bin if using venv)
# The command runs the Flask development server, binding it to all interfaces (0.0.0.0) on port 5000.
# THIS IS NOT RECOMMENDED FOR PRODUCTION.
ExecStart=/path/to/your/flask run --host=0.0.0.0 --port=5000 # Corrected 'flash' to 'flask'

Restart=always

[Install]
WantedBy=multi-user.target
